# salt/files/rsync-backups.timer
[Unit]
Description=rsync backups from primary

[Timer]
OnCalendar=05:30
Persistent=true
Unit=rsync-backups.service

[Install]
WantedBy=timers.target
