# salt/files/rsync-backups.service
[Unit]
Description=rsync backups

[Service]
Type=oneshot
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/local/sbin/rsync-backups.sh
